# src/jdvm-docker/src/main/docker/etc/jdvm-systemd/service/jdvm-apt-update.service
[Unit]
Description=Update list of available packages
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
Environment=DEBIAN_FRONTEND=noninteractive
ExecStart=/usr/bin/apt update
